# systemd/carely-kiosk.service
[Unit]
Description=Carely UI Kiosk (Chromium)
After=carely.service graphical.target
Requires=carely.service

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStartPre=/bin/sleep 3
ExecStart=/usr/bin/chromium-browser \
    --kiosk \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --noerrdialogs \
    --no-first-run \
    --disable-translate \
    --disable-features=TranslateUI \
    http://localhost:8080
Restart=on-failure
RestartSec=5
User=carely

[Install]
WantedBy=graphical.target
