# systemd/carely-audio.service
[Unit]
Description=Carely Audio Pipeline (C++)
After=sound.target
Wants=sound.target

[Service]
Type=simple
ExecStart=/opt/assistant/audio_pipeline/build/carely-audio
Restart=on-failure
RestartSec=2
StartLimitBurst=5
StartLimitIntervalSec=300
Environment=XDG_RUNTIME_DIR=/run/user/1000
EnvironmentFile=/opt/assistant/.env
User=hp
Group=audio
RuntimeDirectory=carely
RuntimeDirectoryMode=0755
Nice=-10
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50
LimitRTPRIO=50
LimitMEMLOCK=infinity
StandardOutput=journal
StandardError=journal
SyslogIdentifier=carely-audio

[Install]
WantedBy=multi-user.target
